# cron/com.viralcommand.daily-discover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.viralcommand.daily-discover</string>
	<key>Program</key>
	<string>__PIPELINE_DIR__/scripts/daily-discover.sh</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>6</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>__PIPELINE_DIR__/logs/daily-discover.log</string>
	<key>StandardErrorPath</key>
	<string>__PIPELINE_DIR__/logs/daily-discover.log</string>
	<key>WorkingDirectory</key>
	<string>__PIPELINE_DIR__</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
	</dict>
</dict>
</plist>
